# infra/systemd/galaxy-os.target.example
[Unit]
Description=Galaxy OS production services
Requires=galaxy-os-api.service galaxy-os-web.service
After=galaxy-os-api.service galaxy-os-web.service

[Install]
WantedBy=default.target
